# deploy/systemd/whatsapp-bridge.service
[Unit]
Description=WhatsApp bridge (whatsmeow session + loopback REST API)
After=network-online.target
Wants=network-online.target
# Three failures in fifteen minutes means the problem is not transient. These
# are [Unit] directives (systemd >= 229); the previous 5s RestartSec was also
# wider than the default 10s window, so no start limit ever engaged.
StartLimitIntervalSec=900
StartLimitBurst=3

[Service]
User=whatsapp
Group=whatsapp
WorkingDirectory=/opt/whatsapp_mcp/bridge
# store/ (whatsapp.db session + messages.db archive + .bridge-token) lives
# under WorkingDirectory on the VM's persistent disk, so the pairing survives
# restarts and redeploys.
#
# Optional (leading '-'): device name and cold-send limits. See .env.example.
# WHATSAPP_ALLOW_PAIRING belongs in the interactive pairing command, not here —
# the service must never be able to start a pairing flow on its own.
EnvironmentFile=-/etc/whatsapp-bridge.env
ExecStart=/opt/whatsapp_mcp/bridge/whatsapp-bridge

# on-failure, not always. The bridge exits 0 when WhatsApp has told it to stop
# (ban, logout, outdated client) and records why in store/.halted; restarting
# into that is what turns a temporary block into a permanent one. Genuine
# crashes still exit non-zero and are retried.
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
